# system/crontab
# shellcheck shell=bash disable=SC2034,SC2035
SHELL=/bin/bash
CRON_TZ=America/Los_Angeles


0 */6 * * *     docker exec kopia -t /usr/local/bin/kopia-backup
